# solar-ac.service
# Solar AC Controller - Serviciu Systemd
# Copiază în: /etc/systemd/system/solar-ac.service
# Activează cu:
#   sudo cp solar-ac.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now solar-ac.service
# Vezi loguri: sudo journalctl -u solar-ac.service -f

[Unit]
Description=Solar AC Controller - Automatizare AC pe baza energiei solare
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=anne
WorkingDirectory=/home/anne/solar-ac-controller
ExecStart=/home/anne/solar-ac-controller/venv/bin/python main.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
